# ec2-docker/.config/super_uwsgi.conf
[program:uwsgi]
command = uwsgi -i /srv/ec2-docker/.config/uwsgi.ini