# eb-docker/.config/super_uwsgi.conf
[program:uwsgi]
command = uwsgi -i /srv/eb-docker/.config/uwsgi.ini